# deploy/vibecheck.service
[Unit]
Description=VibeCheck API (Docker Compose)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/opt/vibecheck
ExecStartPre=/usr/bin/docker compose -f /opt/vibecheck/docker-compose.prod.yml --env-file /opt/vibecheck/.env.prod pull
ExecStart=/usr/bin/docker compose -f /opt/vibecheck/docker-compose.prod.yml --env-file /opt/vibecheck/.env.prod up -d
ExecStop=/usr/bin/docker compose -f /opt/vibecheck/docker-compose.prod.yml --env-file /opt/vibecheck/.env.prod down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
